# server/data/groupgd.service.in
[Service]
ExecStart=@bindir@/groupgd

[Install]
WantedBy=multi-user.target
Also=groupgd.socket
