# server/data/groupgd.service
[Service]
ExecStart=/usr/local/bin/groupgd

[Install]
WantedBy=multi-user.target
Also=groupgd.socket
